# deploy/chat-app-moderation.service
[Unit]
Description=chat-app AI moderation patrol
After=network-online.target chat-app.service
Wants=network-online.target

[Service]
Type=oneshot
User=shohei
Group=shohei
SupplementaryGroups=chat-app
WorkingDirectory=/opt/chat-app/current
EnvironmentFile=/etc/chat-app/chat-app.env
Environment=HOME=/home/shohei
Environment=PATH=/home/shohei/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/node src/moderation-worker.js
UMask=0007

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/chat-app /home/shohei/.claude
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
